# legacy/storyboard/com.opero.storyboard-autopush.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.opero.storyboard-autopush</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/oskarpeterson/Documents/AI/AGENCY /storyboard/auto-sync.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/oskarpeterson/Documents/AI/AGENCY /storyboard/auto-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/oskarpeterson/Documents/AI/AGENCY /storyboard/auto-sync.log</string>
</dict>
</plist>
